# scripts/systemd/etc/systemd/system/IpmiMgr@.service.in
[Unit]
Description=MICROTCA CRATE IPMI MANAGER IOC %p%I
After=rc-local.service
Wants=rc-local.service
After=network-online.target
Wants=network-online.target

[Service]
# Source environment
EnvironmentFile=/etc/sysconfig/ipmi-mgr-epics-ioc
EnvironmentFile=/etc/sysconfig/ipmi-mgr-epics-ioc-mapping
# Execute pre with root
PermissionsStartOnly=true
Environment=IPMIMGR_INSTANCE=%i
ExecStartPre=/bin/mkdir -p /var/log/procServ/%p%i
ExecStartPre=/bin/mkdir -p /var/run/procServ/%p%i
WorkingDirectory=<INSTALL_PREFIX>/<IOC_NAME>/iocBoot/iocIpmiMgr
# Run procServ with IOC
ExecStart="/usr/local/bin/runIpmiMgr.sh"

[Install]
WantedBy=multi-user.target
